# .config/profile.d/21-host-vars.sh
#!/bin/sh

export BROWSER='firefox'
export XTERM='xfce4-terminal'
