# cubicle/configs/.config/profile.d/20-configs-vars.sh
#!/bin/sh

export LC_COLLATE='C'
export PAGER='less --chop-long-lines --jump-target=5 --ignore-case --no-init --RAW-CONTROL-CHARS --quit-if-one-screen --shift=.2'
